# SYN-ISO-PROFILE/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
# synstigator replaces root as the live-session account: wheel (for doas)
# + a known password (so SSH works from boot, no `passwd root` needed
# first) + NOPASSWD via doas so the install scripts don't prompt. wheel's
# GID is allocated at boot by systemd-sysusers (not fixed in any package),
# so this can't be a static /etc/passwd+group edit — it has to run here,
# each boot, before the login prompt appears.
ExecStartPre=/bin/sh -c 'id -u synstigator >/dev/null 2>&1 || { useradd -m -G wheel -s /usr/bin/zsh synstigator; touch /home/synstigator/.zshrc; chown synstigator:synstigator /home/synstigator/.zshrc; }; echo synstigator:synstigator | chpasswd; grep -q "^permit nopass synstigator$" /etc/doas.conf 2>/dev/null || echo "permit nopass synstigator" >> /etc/doas.conf; chmod 600 /etc/doas.conf'
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin synstigator - $TERM
